# .ref-Bell-32V/usr/lib/crontab
20 1 * * * /usr/bin/calendar -
0,5,10,15,20,25,30,35,40,45,50,55 * * * *	/usr/lib/atrun
