# usr/src/etc/rc.local
#
# site-specific startup actions, daemons
#
#	@(#)rc.local	8.1 (Berkeley) 6/10/93
#

T=/tmp/_motd
rm -f $T
strings /vmunix | grep UNIX > $T
echo "" >> $T
sed '1,/^$/d' < /etc/motd >> $T
cp $T /etc/motd
chmod 666 /etc/motd
rm -f $T

echo -n 'starting local daemons:'

if [ -f /usr/contrib/lib/jove/recover ]; then
	echo -n ' jove'
	(/usr/contrib/lib/jove/recover -syscrash)
fi

# Kerberos runs ONLY on the Kerberos server machine
#echo -n ' kerberos';	kerberos >> /var/log/kerberos.log &

echo '.'
