# .ref-BSD-4_3_Reno/usr/src/etc/rc.local
#
# site-specific startup actions, daemons
#
#	@(#)rc.local	5.3 (Berkeley) 6/30/90
#

strings /vmunix | grep UNIX >/tmp/t1
tail +2 /etc/motd >>/tmp/t1
mv /tmp/t1 /etc/motd
chmod 666 /etc/motd

echo -n 'starting local daemons:'

# Kerberos runs ONLY on the Kerberos server machine
#echo -n ' kerberos';	kerberos >> /var/log/kerberos.log &

echo '.'
