# .ref-83d6686a4b21d55e8eb1fff796741ef79d4ae20a/usr/src/etc/rc.local
#
# site-specific startup actions, daemons
#
#	@(#)rc.local	5.1 (Berkeley) %G%
#

strings /vmunix | grep UNIX >/tmp/t1
tail +2 /etc/motd >>/tmp/t1
mv /tmp/t1 /etc/motd
chmod 666 /etc/motd

echo -n			starting local daemons:			>/dev/console
#if [ -f /usr/sbin/timed ]; then
#	timed -M &			echo -n ' timed'	>/dev/console
#fi
if [ -f /usr/lib/sendmail ]; then
	(cd /usr/spool/mqueue; rm -f lf*)
	/usr/lib/sendmail -bd -q30m; echo -n ' sendmail'	>/dev/console
fi
				echo '.'			>/dev/console
