# usr/src/contrib/mh-6.8/support/bboards/crontab
00 5 * * * /usr/contrib/mh-6.8/lib/conflict -mail PostMaster
05 5 * * * su bboards < /usr/spool/bboards/etc/bboards.daily
10 5 * * 0 su bboards < /usr/spool/bboards/etc/bboards.weekly
